# InPi/rc.local
#!/bin/bash
sudo -i -u ubuntu /home/ubuntu/start-test-harness.sh