# artifacts/rc.local
#!/bin/bash
sudo -i -u ubuntu /home/ubuntu/start_test_harness.sh
